# .local/share/crontab
*/2 * * * * $HOME/.local/bin/ufiledb
*/2 * * * * $HOME/.local/bin/udirdb
*/15 * * * * /usr/local/bin/newsboat -x reload
