# .local/share/crontab
*/15 * * * * /usr/local/bin/newsboat -x reload
